# docker/assets/supervisord.conf
[supervisord]
nodaemon=true

[program:lumino]
command=/root/lumino/startLumino
autostart=true
user=root
startsecs=0
startretries=0
exitcodes=0